# curdow.timer
[Unit]
Description=Run curdow service periodically

[Timer]
OnCalendar=*-*-* 15:40:00
Persistent=true
Unit=curdow.service
